# crontab
0 12 * * * root /vancouver_updater/run.sh >> /var/log/cron.log
